# pinger.timer
[Unit]
Description=Run Linux Status Pinger every 12 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
