# pinger.timer
[Unit]
Description=Run Linux Status Pinger hourly (sends only at configured REPORT_HOURS)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
